# dist/launchd/dev.terminal-hub.tmux.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.terminal-hub.tmux</string>
  <key>ProgramArguments</key>
  <array>
    <string>__TMUX_BIN__</string>
    <string>-L</string><string>terminal-hub</string>
    <string>new-session</string><string>-d</string>
    <string>-s</string><string>_boot</string>
    <string>sleep</string><string>infinity</string>
  </array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/terminal-hub-tmux.log</string>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/terminal-hub-tmux.log</string>
</dict>
</plist>
