# dist/deb/systemd/tmux-server.service
[Unit]
Description=tmux server backing terminal-hub
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/tmux -L terminal-hub new-session -d -s _boot 'sleep infinity'
ExecStartPost=/usr/bin/tmux -L terminal-hub set-option -g remain-on-exit on
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
