# dist/systemd/terminal-hub.service
[Unit]
Description=terminal-hub web server
After=tmux-server.service
Wants=tmux-server.service

[Service]
Type=simple
ExecStart=__INSTALL_PREFIX__/bin/terminal-hub
Environment="TERMINAL_HUB_BIND=127.0.0.1:5999"
Environment="TERMINAL_HUB_PUBLIC_URL=https://localhost:5999/"
Environment="TERMINAL_HUB_CONFIG_DIR=__HOME__/.config/terminal-hub"
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
